# patchdb/logrotate/logrotate.timer
[Unit]
Description=Rotate log files daily

[Timer]
OnCalendar=Mon..Sun 22:30:00
Persistent=true

[Install]
WantedBy=timers.target
